# overlay/etc/rc.local
#!/bin/bash

#Set up the live environment.

chmod 0440 /etc/sudoers 
mkdir -p /home/synth 
useradd -u 1000 -g users -G storage,optical,audio,video,network,wheel,power,hal,lp -d /home/synth synth
passwd -d synth 
chown synth /home/synth
su -c "cp -r /etc/skel/.[a-zA-Z0-9]* /home/synth/" synth

sleep 10

#TODO: Put something in here to set the first soundcard to 80% volume or so
clear
